# system-config/polargraph-mqtt.service
[Unit]
Description=Polargraph MQTT Subscriber
After=network.target polargraph.service
Requires=polargraph.service

[Service]
Type=simple
User=monalisa
WorkingDirectory=/home/monalisa/polargraph/src
ExecStart=/home/monalisa/polargraph/venv/bin/python mqtt_subscriber.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

